# scripts/breadnotes.service
[Unit]
Description=BreadNotes FastAPI Application
After=network.target

[Service]
Type=exec
User=root
Group=docker
WorkingDirectory=/opt/breadnotes
EnvironmentFile=/opt/breadnotes/.env
ExecStartPre=/usr/bin/docker stop breadnotes-api || true
ExecStartPre=/usr/bin/docker rm breadnotes-api || true
ExecStart=/usr/bin/docker run --name breadnotes-api --rm --env-file /opt/breadnotes/.env -p 8000:8000 breadnotes-api
ExecStop=/usr/bin/docker stop breadnotes-api
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target