# scripts/breadnotes.service
[Unit]
Description=BreadNotes FastAPI Application
After=network.target

[Service]
Type=exec
User=breadnotes
Group=breadnotes
WorkingDirectory=/opt/breadnotes/api
Environment=PATH=/opt/breadnotes/.venv/bin
EnvironmentFile=/opt/breadnotes/.env
ExecStart=/opt/breadnotes/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/breadnotes /var/log/breadnotes

[Install]
WantedBy=multi-user.target